# .config/systemd/user/playerctld.service
[Unit]
Description=Start playerctld

[Service]
Type=simple
ExecStart=/usr/bin/playerctld

[Install]
WantedBy=default.target
